# Raspberrypi/Services/ServiceConfigs/bluetooth-serialkey.service
[Unit]
Description="Bluetooth agent"
After=multi-user.target
#PartOf=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/share/doc/bluez-test-scripts/examples/simple-agent 
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
